# server-monitor-bot.service
[Unit]
Description=Telegram Server Monitor Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/server-monitor-bot
ExecStart=/usr/bin/python3 /opt/server-monitor-bot/server_monitor_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
